# supervisord.conf
[supervisord]
nodaemon=true

[program:mongodb]
command=/usr/bin/mongod --dbpath /mongodb/dbpath/

[program:parse]
command=/usr/local/bin/node /parse/index.js
